# docker/wis2box/wis2box.cron
0 0 * * * wis2box wis2box data clean --days=$WIS2BOX_DATA_RETENTION_DAYS > /proc/1/fd/1 2>/proc/1/fd/2
0 1 * * * wis2box wis2box data archive > /proc/1/fd/1 2>/proc/1/fd/2
*/10 * * * * wis2box echo "wis2box.cron is alive" > /proc/1/fd/1 2>/proc/1/fd/2
